# neptcheck.service
[Service]
ExecStart=/bin/bash /home/core/neptcheck.sh
Type=oneshot
